# operations/supervisord.d/laravel-globus-worker.ini
[program:laravel-globus-worker]
process_name = %(program_name)s_%(process_num)02d
command = php /var/www/html/materialscommons/artisan queue:work --queue=globus --sleep=3 --tries=3 --timeout=3600 --delay=10
autostart = true
autorestart = true
user = nginx
numprocs = 4
redirect_stderr = true
stdout_logfile = /usr/local/miserver/logs/materialscommons/globus-worker.log
