# operations/supervisord.d/laravel-spreadsheets-worker.ini
[program:laravel-spreadsheets-worker]
process_name = %(program_name)s_%(process_num)02d
command = php /var/www/html/materialscommons/artisan queue:work --queue=spreadsheets --sleep=3 --tries=3 --timeout=900 --delay=10
autostart = true
autorestart = true
user = nginx
numprocs = 4
redirect_stderr = true
stdout_logfile = /usr/local/miserver/logs/materialscommons/spreadsheet-worker.log
